# cuddy.service
[Unit]
Description=cuddy

[Service]
ExecStart=/usr/local/bin/cuddy
Restart=always

[Install]
WantedBy=multi-user.target
